# bambuboard.service
[Unit]
Description=BambuBoard
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/bambuboard
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=3
User=bambuboard
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
